# ops/systemd/ai-ops-orb-daily.timer
[Unit]
Description=ai-ops-runner daily ORB check timer

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
