# ops/systemd/openclaw-novnc-shm-guard.timer
[Unit]
Description=OpenClaw noVNC SHM guard timer (every 5 minutes)

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
