# ops/systemd/openclaw-doctor.service
[Unit]
Description=OpenClaw infrastructure health + audit check
After=ai-ops-runner.service
Requires=ai-ops-runner.service

[Service]
Type=oneshot
WorkingDirectory=/opt/ai-ops-runner
ExecStart=/bin/bash /opt/ai-ops-runner/ops/openclaw_doctor.sh
TimeoutStartSec=120
