# ops/systemd/openclaw-novnc.service
[Unit]
Description=OpenClaw noVNC (Xvfb + x11vnc + websockify) for Kajabi interactive capture
After=network.target

[Service]
Type=simple
# Explicit DISPLAY and ports (overridden by next.env when capture writes it)
Environment=OPENCLAW_NOVNC_ARTIFACT_DIR=/run/openclaw-novnc
Environment=OPENCLAW_NOVNC_PORT=6080
Environment=OPENCLAW_NOVNC_DISPLAY=:99
Environment=OPENCLAW_NOVNC_VNC_PORT=5900
# Params from next.env (written by kajabi_capture_interactive.py before start)
EnvironmentFile=-/run/openclaw-novnc/next.env

# Pre-start: remove stale X lock for configured DISPLAY only (safe, minimal)
ExecStartPre=/bin/bash -c 'd=${OPENCLAW_NOVNC_DISPLAY#:}; d=${d:-99}; L=/tmp/.X${d}-lock; [ -f "$L" ] && p=$(cat "$L" 2>/dev/null) && [ -n "$p" ] && ! kill -0 "$p" 2>/dev/null && rm -f "$L" || true'

ExecStart=/opt/ai-ops-runner/ops/novnc_supervisor.sh
WorkingDirectory=/opt/ai-ops-runner
Restart=always
RestartSec=2
TimeoutStartSec=15
TimeoutStopSec=15

# Clean shutdown: supervisor traps TERM and kills children
KillMode=mixed
KillSignal=SIGTERM

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-novnc

[Install]
# Not enabled by default; started on demand by kajabi_capture_interactive
WantedBy=multi-user.target
