# ops/systemd/openclaw-rootd-watchdog.timer
[Unit]
Description=OpenClaw rootd watchdog (every 60s)

[Timer]
OnBootSec=60
OnUnitActiveSec=60
Persistent=true

[Install]
WantedBy=timers.target
